# config/cron/crontab
# Contenedor Crontab
# Formato: minuto hora día mes día_semana comando
#
# Variables de entorno necesarias para cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app
PYTHONUNBUFFERED=1

# Redirigir toda la salida al log
CRON_LOG=/var/log/cron/cron.log

# =============================================================================
# TAREAS PROGRAMADAS
# =============================================================================

# iptvListWatcher - Descargar y verificar lista IPTV cada dia a las 05:00 y a las 14:00
0 5,14 * * * cd /app && python -m iptvListWatcher download >> /app/logs/iptvListWatcher.log 2>&1

# iptvListWatcher - Limpiar backups antiguos (mantener últimos 5) cada semana, los miercoles a las 7:00 AM
0 3 * * 3 cd /app && python -m iptvListWatcher clean --keep 5 >> /app/logs/iptvListWatcher.log 2>&1

# =============================================================================
# TAREAS DE MANTENIMIENTO
# =============================================================================

# Limpiar logs antiguos de más de 7 días (cada lunes a las 5:00 AM)
0 5 * * 1 find /app/logs -name "*.log" -mtime +7 -delete >> $CRON_LOG 2>&1

# Rotar el log de cron si supera 10MB con timestamp (diario a las 4:00 AM)
0 4 * * * [ $(stat -c\%s $CRON_LOG 2>/dev/null) -gt 10485760 ] && mv $CRON_LOG $CRON_LOG.old.$(date +\%Y\%m\%d-\%H\%M\%S) && touch $CRON_LOG 2>&1

